# dotfiles/.vim/template/crontab/base-main.crontab
*/1 * * * * ruby ~/Dropbox/dev/prj/bottalo/bottalo.rb > /dev/null 2>&1
